# .config/systemd/user/xbindkeys.service
[Unit]
Description = xbindkeys

[Service]
ExecStart = /usr/bin/xbindkeys --nodaemon

[Install]
WantedBy = default.target